# scripts/touchboost_openrc
#!/sbin/openrc-run

# Prefix: PREFIX

command="PREFIX/bin/touchboost"
command_args="`PREFIX/bin/touchboost_detect`"
pidfile="/var/run/touchboost.pid"
name="Touchboost"

depend() {
	need udev
}
